# deploy/systemd/dona-cron.service
[Unit]
Description=Dona Cron service
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=cronapp
Group=cronapp
WorkingDirectory=/var/www/cron
EnvironmentFile=/etc/dona-cron/cron.env
ExecStart=/usr/bin/java -jar /var/www/cron/target/cron-service-0.0.1-SNAPSHOT.jar
Restart=always
RestartSec=5
SuccessExitStatus=143

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/www/cron/target

[Install]
WantedBy=multi-user.target
